# templates/default/delayed_job.conf.erb
description "<%= @application['id'] %> delayed_job"

start on (local-filesystems and runlevel [2345])
stop on runlevel [016]

respawn
respawn limit 10 5

pre-start script
  mkdir -p /var/run/<%= @application['id'] %>
  chown -R <%=@deployment[:user]%>:<%=@deployment[:user]%> /var/run/<%= @application['id'] %>
end script

exec start-stop-daemon \
  --start \
  --pidfile /var/run/<%= @application['id'] %>/<%= @application['id'] %>.pid \
  --chuid <%=@deployment[:user]%>:<%=@deployment[:user]%> \
  --chdir <%= @deployment[:deploy_to] %>/current \
  --exec /usr/bin/env RAILS_ENV=<%= @application['rails_env'] || node.chef_environment %> /usr/local/bin/bundle exec rake jobs:work
