# templates/default/sidekiq_workers.conf.erb
# from:  https://github.com/mperham/sidekiq/blob/master/examples/upstart/manage-one/workers.conf

# manage a set of Sidekiqs

# This config should work with Ubuntu 12.04+.  It
# allows you to manage multiple Sidekiq instances with
# Upstart, Ubuntu's native service management tool.
#
# See sidekiq.conf for how to manage a single Sidekiq instance.
#
# Use "stop sidekiq_workers_app" to stop all Sidekiq instances.
# Use "start sidekiq_workers_app" to start all instances.
# Use "restart sidekiq_workers_app" to restart all instances.
# Crazy, right?
#

description "sidekiq workers management for '<%= @application['id'] %>'"

# This starts upon bootup and stops on shutdown
start on runlevel [2345]
stop on runlevel [06]

# Set this to the number of Sidekiq processes you want
# to run on this machine
env NUM_WORKERS=<%= @application['sidekiq']['num_workers'] || 1 %>

pre-start script
  for i in `seq 0 $((${NUM_WORKERS} - 1))`
  do
    start sidekiq_<%= @application['id'] %> index=$i
  done
end script